# supplemental/debian/vigil-agent.service
[Unit]
Description=Vigil Agent Service
Wants=network-online.target
After=network-online.target

[Service]
EnvironmentFile=/etc/vigil-agent.conf
ExecStart=/usr/bin/vigil-agent
User=vigil
Group=vigil
Restart=on-failure
RestartSec=5
StateDirectory=vigil-agent

# Security/sandboxing settings
NoNewPrivileges=yes
PrivateTmp=yes
KeyringMode=private
LockPersonality=yes
ProtectClock=yes
ProtectHome=read-only
ProtectHostname=yes
ProtectKernelLogs=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
ProtectSystem=strict
RemoveIPC=yes
RestrictSUIDSGID=true

[Install]
WantedBy=multi-user.target
